# services/docker/supervisord.conf
[program:docker]
autostart=true
command=/bin/bash -c "exec /usr/bin/dockerd --raw-logs"
stdout_logfile=/var/log/cavirin/docker.log
redirect_stderr=true
priority=100
