# services/nginx/supervisord.conf
[program:nginx]
command=/bin/bash -c 'exec /usr/sbin/nginx -g "daemon off;"'
stdout_logfile=/var/log/cavirin/nginx.log
redirect_stderr=true
priority=150
autorestart=true
