# services/tomcat8/supervisord.conf
[program:tomcat8]
command=/usr/share/tomcat8/bin/tomcatd.sh
stdout_logfile=/var/log/cavirin/tomcat.log
redirect_stderr=true
priority=100
