# services/zookeeper/supervisord.conf
[program:zookeeper]
command=/bin/bash -c "/usr/share/zookeeper/bin/zkServer.sh start-foreground"
stdout_logfile=/var/log/zookeeper/zookeeper.log
redirect_stderr=true
priority=100
