# services/sshd/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/cavirin/sshd.log
redirect_stderr=true
priority=50
