# services/vault/supervisord.conf
[program:vault-server]
user=ubuntu
command=/opt/cavirin/vaultdir/vaultd.sh
stdout_logfile=/var/log/cavirin/vault-sv.log
redirect_stderr=true
autorestart=true
priority=100
